# openwrt/files/etc/crontabs/root
*/5 * * * * bash /root/node/wan_checker.sh

